# detection-rules/alerts.conf
# Splunk Detection Rules — SIEM Homelab
# deploy via: python3 deploy_alerts.py
# or manually import this file into Splunk

# ============================================================
# RULE 1 — SSH Brute Force Detection (Linux)
# Triggers when failed SSH login attempts exceed 5 from any host
# MITRE ATT&CK: T1110.001 — Brute Force: Password Guessing
# ============================================================
[HOMELAB - SSH Brute Force Detected]
search = index=main sourcetype=linux_secure "Failed password" | stats count by host | where count > 5
dispatch.earliest_time = -5m
dispatch.latest_time = now
cron_schedule = */5 * * * *
enableSched = 1
alert.severity = 2
alert.track = 1
description = Detects SSH brute force attacks by counting failed password attempts. Triggers when count exceeds 5 from a single host within 5 minutes.

# ============================================================
# RULE 2 — Failed Windows Login Detection (Event ID 4625)
# Triggers when failed Windows logon attempts exceed 3
# MITRE ATT&CK: T1110.001 — Brute Force: Password Guessing
# ============================================================
[HOMELAB - Brute Force Failed Windows Logins]
search = index=main sourcetype="WinEventLog:Security" EventCode=4625 | stats count by Account_Name, host | where count > 3
dispatch.earliest_time = -5m
dispatch.latest_time = now
cron_schedule = */5 * * * *
enableSched = 1
alert.severity = 2
alert.track = 1
description = Detects repeated failed Windows login attempts. Triggers when a single account exceeds 3 failed logons within 5 minutes. Maps to credential access techniques.

# ============================================================
# RULE 3 — New User Account Creation (Event ID 4720)
# Triggers on any new Windows user account creation
# MITRE ATT&CK: T1136 — Create Account
# ============================================================
[HOMELAB - New User Account Created]
search = index=main sourcetype="WinEventLog:Security" EventCode=4720
dispatch.earliest_time = -5m
dispatch.latest_time = now
cron_schedule = */5 * * * *
enableSched = 1
alert.severity = 3
alert.track = 1
description = Detects creation of new Windows user accounts (Event ID 4720). Any new account creation in a monitored environment should be investigated to confirm it is authorized.

# ============================================================
# RULE 4 — Password Spraying Detection (Linux)
# Triggers when one source IP targets more than 3 different usernames
# MITRE ATT&CK: T1110.003 — Password Spraying
# ============================================================
[HOMELAB - Password Spraying Detected]
search = index=main sourcetype=linux_secure "Failed password" | rex field=_raw "for (?:invalid user )?(?P<username>\S+) from (?P<src_ip>\S+)" | bucket _time span=5m | stats dc(username) as unique_users count as attempts by src_ip _time | where unique_users > 3 AND attempts > 10
dispatch.earliest_time = -5m
dispatch.latest_time = now
cron_schedule = */5 * * * *
enableSched = 1
alert.severity = 2
alert.track = 1
description = Detects password spraying — one source IP targeting multiple usernames with the same password. Key difference from brute force: dc(username) > 3 from the same source IP in 5 minutes.

# ============================================================
# RULE 5 — New Admin Account Persistence (Event ID 4720)
# Triggers on new account creation with enriched context
# MITRE ATT&CK: T1078 — Valid Accounts
# ============================================================
[HOMELAB - New Admin Account Persistence]
search = index=main sourcetype="WinEventLog:Security" EventCode=4720 | table _time src_user user dest
dispatch.earliest_time = -5m
dispatch.latest_time = now
cron_schedule = */5 * * * *
enableSched = 1
alert.severity = 3
alert.track = 1
description = Detects new local account creation with full context — who created it and on which host. Attackers create accounts for persistence after gaining initial access.

# ============================================================
# RULE 6 — Scheduled Task Created (Event ID 4698)
# Triggers on any scheduled task creation
# MITRE ATT&CK: T1053.005 — Scheduled Task/Job
# ============================================================
[HOMELAB - Scheduled Task Created]
search = index=main sourcetype="WinEventLog:Security" EventCode=4698 | table _time src_user TaskName TaskContent
dispatch.earliest_time = -5m
dispatch.latest_time = now
cron_schedule = */5 * * * *
enableSched = 1
alert.severity = 3
alert.track = 1
description = Detects scheduled task creation (Event ID 4698). Attackers use scheduled tasks for persistence — code runs automatically on reboot or on a timer without requiring re-entry.

# ============================================================
# RULE 7 — Windows Event Logs Cleared (Event ID 1102 / 104)
# Triggers when security or system logs are cleared
# MITRE ATT&CK: T1070.001 — Indicator Removal: Clear Windows Event Logs
# ============================================================
[HOMELAB - Windows Event Logs Cleared]
search = index=main sourcetype="WinEventLog:Security" (EventCode=1102 OR EventCode=104) | table _time src_user Message
dispatch.earliest_time = -5m
dispatch.latest_time = now
cron_schedule = */5 * * * *
enableSched = 1
alert.severity = 2
alert.track = 1
description = Detects Windows event log clearing (1102 = Security log cleared, 104 = System log cleared). Attackers clear logs to remove evidence of their activity. Almost zero false positives in production.
